# .provision/supervisord_vault.conf
[program:vault]
command=/usr/sbin/vault server -config=/etc/vault/server.conf
priority=1