# .docker/redcap/config/supervisor/celery.ini
[program:celery]
command=celery -A redcap worker -l info
directory=/redcap

