# .docker/redcap/config/supervisor/gunicorn.ini
[program:gunicorn]
command=gunicorn -b 0.0.0.0:8000 redcap.wsgi 
directory=/redcap
priority=1
